# files/system/etc/systemd/user/neovim-setup.service
[Unit]
Description=First-login Neovim setup via bob
ConditionPathExists=!/home/%u/.config/neovim-setup-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-neovim.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target
